# csbridge.service
[Unit]
Description=csbridge

[Service]
ExecStart=/usr/bin/csbridge.py

[Install]
WantedBy=multi-user.target