# deb/etc/systemd/system/csbridge.service
[Unit]
Description=csbridge service

[Service]
User=csbridge
ExecStart=/usr/bin/csbridge.py

[Install]
WantedBy=multi-user.target